# server/templates/container-Jellyfin.service.j2
# container-Jellyfin.service
# autogenerated by Podman 3.4.4
# Sun Jan  2 16:51:43 CET 2022

[Unit]
Description=Podman container-Jellyfin.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman container run --cidfile=%t/%n.ctr-id --cgroups=no-conmon --rm --sdnotify=conmon --replace --name Jellyfin --label io.containers.autoupdate=registry --publish 8096:8096/tcp --userns keep-id --user 1000:1000 --volume {{  jellyfin_cache }}:/cache:Z --volume {{ jellyfin_config }}:/config:Z --volume {{ jellyfin_media }}:/media:ro,z --detach=True docker.io/jellyfin/jellyfin:latest
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm -f --ignore --cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all
User=martin

[Install]
WantedBy=default.target
